# src/k8s/smc-project/image/crontab
#!/bin/sh
/usr/bin/crontab "$@"
/usr/bin/crontab -l > $HOME/.crontab_smc